# nohup/supervisord.conf
[supervisord]
nodaemon=true
user=appuser

[program:test]
command=nohup /bin/bash /home/test.sh &
